# etc/flux-data.service.in
[Unit]
Description=Flux data service
BindsTo=flux.service
After=flux.service

[Service]
TimeoutStopSec=90
KillMode=mixed
ExecStart=@X_BINDIR@/flux data-service
SyslogIdentifier=flux-data
Restart=always
RestartSec=10s

User=flux
Group=flux

[Install]
WantedBy=flux.service
